# supervisord.conf
[program:myapp]
command=python -m androidmonitor_backend.background
directory=.
autostart=true
autorestart=true