# supervisord.conf
[program:myapp]
command=python -m androidmonitor_backend.backend
directory=.
autostart=true
autorestart=true